# brightness-night.service
[Unit]
Description=Set screen brightness to 10% for nighttime

[Service]
Type=oneshot
ExecStart=/usr/bin/brightnessctl set 10%
User=root

